# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=riff

[program:riff]
command=java -jar /home/riff/riff-web.jar

[inet_http_server]
port=*:9001
